# .service
[Unit]
Description=Gunicorn instance to serve {{ app_name }}
After=network.target

[Service]
User=azure-admin
Group=www-data
WorkingDirectory=/home/azure-admin/flask-ansible-example
Environment="PATH=/home/azure-admin/flask-ansible-example/env/bin"
ExecStart=/home/azure-admin/flask-ansible-example/env/bin/gunicorn --workers 3 --bind unix:/home/azure-admin/flask-ansible-example.sock -m 007 wsgi:app

[Install]
WantedBy=multi-user.target
